# web13_pycharm_github.conf
[program:web13_pycharm_github]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/web13_pycharm_github
autostart=true
autorestart=true
