# web13_pycharm_github.conf
[program:web13_pycharm_github]
directory=/root/web13_pycharm_github
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
autostart=true
autorestart=true
